# recipes-core/bit/files/bit.service
[Unit]
Description=Bit Service

[Service]
ExecStart=/etc/TestServer_evb

[Install]
WantedBy=multi-user.target

