# recipes-core/bit/files/bit.service
[Unit]
Description=BIT Service

[Service]
ExecStart=/etc/TestServer_evb

[Install]
WantedBy=multi-user.target
