# recipes-core/bit/files/bit.service
[Unit]
Description=Bit Service

[Service]
ExecStart=/usr/bin/bit

[Install]
WantedBy=multi-user.target
