# recipes-core/bit/files/bit.service
[Unit]
Description=BIT Service

[Service]
ExecStart=/usr/bin/bit.sh

[Install]
WantedBy=multi-user.target

