# recipes-core/maintenance/files/maintenance.service
[Unit]
Description=Maintenance Service

[Service]
ExecStart=/usr/bin/maintenance.sh

[Install]
WantedBy=multi-user.target

